# rc.local
#!/bin/bash

cd /home/ubuntu/poll
source venv/bin/activate

while true
do
    python poll.py
done